# deploy/systemd/ephyr-broker.service
[Unit]
Description=Ephyr SSH Certificate Broker
Documentation=https://github.com/EphyrAI/Ephyr
After=network.target ephyr-signer.service
Wants=ephyr-signer.service

[Service]
Type=simple
User=ephyr-broker
Group=ephyr-agents

ExecStart=/usr/local/bin/ephyr broker --policy /etc/ephyr/policy.yaml
ExecReload=/bin/kill -HUP $MAINPID

# Port and token configuration (defaults shown).
# Override via systemd drop-in or environment file.
Environment=EPHYR_SIGNER_SOCKET=/run/ephyr/signer.sock
# Environment=EPHYR_MCP_LISTEN=:8554
# Environment=EPHYR_DASHBOARD_LISTEN=:8553
# Environment=EPHYR_DASHBOARD_TOKEN=changeme

StandardOutput=journal
StandardError=journal
SyslogIdentifier=ephyr-broker

Restart=on-failure
RestartSec=5s

# Security hardening.
ProtectSystem=strict
ProtectHome=yes
NoNewPrivileges=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes
RestrictNamespaces=yes

ReadWritePaths=/run/ephyr /var/log/ephyr /var/lib/ephyr

CapabilityBoundingSet=
AmbientCapabilities=
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

[Install]
WantedBy=multi-user.target
