# deploy/systemd/ephyr-signer.service
[Unit]
Description=Ephyr SSH Certificate Signer
Documentation=https://github.com/EphyrAI/Ephyr
After=network.target
Before=ephyr-broker.service

[Service]
Type=simple
User=ephyr-broker
Group=ephyr-agents

ExecStart=/usr/local/bin/ephyr signer --ca-key /etc/ephyr/ca_key --socket /run/ephyr/signer.sock

StandardOutput=journal
StandardError=journal
SyslogIdentifier=ephyr-signer

Restart=on-failure
RestartSec=5s

# Security hardening.
ProtectSystem=strict
ProtectHome=yes
NoNewPrivileges=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes
RestrictNamespaces=yes

ReadOnlyPaths=/etc/ephyr
ReadWritePaths=/run/ephyr

CapabilityBoundingSet=
AmbientCapabilities=
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

[Install]
WantedBy=multi-user.target
